# aptly-webui/systemd/aptly-webui-backend.service
[Unit]
Description=Aptly WebUI Backend
After=network.target aptly-api.service
Wants=aptly-api.service

[Service]
Type=simple
User=aptly
WorkingDirectory=/home/keith/aptlywebui/aptly-webui/backend
Environment=APTLY_API_URL=http://10.0.2.160:5000
Environment=CACHE_DB_PATH=/var/lib/aptly-webui/cache.db
Environment=USE_OPTIMIZED_SYNC=true
Environment=ENABLE_AUTO_SYNC=true
Environment=FLASK_DEBUG=false
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/python3 -m flask --app app run --host=0.0.0.0 --port=5001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
